# deploy/systemd/doc2vault-watch@.service
# systemd-Template fuer die dauerhafte Ordnerueberwachung eines Jobs.
#
# Installation (als root):
#   1. Pfade unten anpassen (ExecStart, User, ggf. Environment).
#   2. Datei nach /etc/systemd/system/doc2vault-watch@.service kopieren.
#   3. systemctl daemon-reload
#   4. systemctl enable --now doc2vault-watch@<job-id>
#      (Job-ID aus `doc2vault-jobs list`, z. B. berichte-2024)
#
# Der Dienst laeuft unter dem angegebenen Nutzer; Jobs/Manifeste liegen in
# dessen Konfigverzeichnis (~/.config/doc2vault bzw.
# DOC2VAULT_HOME).

[Unit]
Description=doc2vault-Ordnerueberwachung (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Anpassen: Pfad zur Installation (venv-Beispiel).
ExecStart=/opt/doc2vault/.venv/bin/doc2vault-jobs watch %i
# Anpassen: Nutzer, dem der Job gehoert.
User=docling
Restart=on-failure
RestartSec=10
Environment=PYTHONUNBUFFERED=1
# Optional: abweichendes Konfigverzeichnis.
# Environment=DOC2VAULT_HOME=/srv/doc2vault

[Install]
WantedBy=multi-user.target
